# debian/mogilefsd.prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/mogilefsd" ] || [ -e "/etc/init/mogilefsd.conf" ]; then
	invoke-rc.d mogilefsd stop || exit $?
fi